# net-misc/apache-apollo/files/apollo
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

BROKER=${SVCNAME#*.}
if [ -z "${BROKER}" ] || [ "${BROKER}" == "apollo" ]; then
    BROKER="default"
fi
PIDFILE="/var/run/apollo.${BROKER}.pid"

depend() {
    need localmount net
    after bootmisc
}

start() {
    ebegin "Starting apollo broker ${BROKER}"
    start-stop-daemon \
        --start --quiet --background \
        --pidfile "${PIDFILE}" --make-pidfile \
        --user apollo --group apollo \
        --exec "/var/lib/apollo/${BROKER}/bin/apollo-broker" \
        -- run
    eend $?
}

stop() {
    ebegin "Stopping apollo broker ${BROKER}"
    start-stop-daemon \
        --stop --quiet --pidfile "${PIDFILE}"
    eend $?
}

