# media-sound/teamspeak-server-bin/files/initd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
	use mysql
}

start() {
	ebegin "Starting TeamSpeak 3 Server"

	[ -d /run/teamspeak3-server ] || mkdir /run/teamspeak3-server
	chown teamspeak3:teamspeak3 /run/teamspeak3-server

	start-stop-daemon --start --background \
		--pidfile "/run/teamspeak3-server/server.pid" --make-pidfile \
		--user "teamspeak3" --chdir "/opt/teamspeak3-server" \
		--exec "/usr/bin/ts3server" -- \
		inifile="/etc/teamspeak3-server/server.conf"

	eend $?
}

stop() {
	ebegin "Stopping TeamSpeak Server"

	start-stop-daemon --stop \
		--pidfile "/run/teamspeak3-server/server.pid"

	eend $?
}

restart() {
	svc_stop
	sleep 3
	svc_start
}
