# media-sound/teamspeak-server-bin/files/initd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
	use mysql
}

start() {
	ebegin "Starting teamspeak3 server"

	[ -d /run/teamspeak3-server ] || mkdir /run/teamspeak3-server
	chown teamspeak3:teamspeak3 /run/teamspeak3-server

	start-stop-daemon --start --background \
		--pidfile "/run/teamspeak3-server/server.pid" --make-pidfile \
		--user "teamspeak3" \
		--exec "/usr/bin/ts3server" -- \
		inifile="/etc/teamspeak3-server/server.conf"

	eend $?


	[ -f /usr/bin/tsdnsserver ] || return


	ebegin "Starting tsdns server"

	[ -f /dev/shm/tsdns_update ] && rm /dev/shm/tsdns_update

	start-stop-daemon --start --background \
		--pidfile "/run/teamspeak3-server/tsdns.pid" --make-pidfile \
		--stdout "/var/log/teamspeak3-server/tsdns.log" --stderr "/var/log/teamspeak3-server/tsdns.err" \
		--user "teamspeak3" \
		--exec "/usr/bin/tsdnsserver"

	eend $?
}

stop() {
	ebegin "Stopping teamspeak3 server"

	start-stop-daemon --stop \
		--pidfile "/run/teamspeak3-server/server.pid"

	eend $?


        [ -f /usr/bin/tsdnsserver ] || return


	ebegin "Stopping tsdns server"

	start-stop-daemon --stop \
		--pidfile "/run/teamspeak3-server/tsdns.pid"

	eend $?
}

restart() {
	svc_stop
	sleep 3
	svc_start
}
