# net-p2p/syncthing/files/syncthing_initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDFILE="/run/syncthing/syncthing.pid"

depend() {
    need localmount net
    after bootmisc
}

start() {
    [ -d /run/syncthing ] || mkdir /run/syncthing
    [ -d /var/log/syncthing ] || mkdir /var/log/syncthing

    chown -R "${SYNCTHING_USER}:${SYNCTHING_GROUP}" /run/syncthing /var/log/syncthing

    ebegin "Starting syncthing"
    start-stop-daemon \
        --start \
        --pidfile "${PIDFILE}" \
        --make-pidfile \
        --stdout "/var/log/syncthing/syncthing.log" \	
        --stderr "/var/log/syncthing/syncthing.err" \
        --user "${SYNCTHING_USER}" \
        --group "${SYNCTHING_GROUP}" \
        --env "STNORESTART=yes" \
        --background \
        --exec /usr/bin/syncthing
    eend $?
}

stop() {
    ebegin "Stopping syncthing"
    start-stop-daemon \
        --stop \
        --pidfile "${PIDFILE}"
    eend $?
}
