# net-p2p/syncthing/files/syncthing_initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $


depend() {
    need localmount net
    after bootmisc
}

start() {
    [ -d /var/log/syncthing ] || mkdir /var/log/syncthing
    chown -R "${SYNCTHING_USER}:${SYNCTHING_GROUP}" /var/log/syncthing

    ebegin "Starting syncthing"
    start-stop-daemon \
        --start \
        --stdout "/var/log/syncthing/syncthing.log" \
        --stderr "/var/log/syncthing/syncthing.err" \
        --user "${SYNCTHING_USER}" \
        --group "${SYNCTHING_GROUP}" \
        --background \
        --exec /usr/bin/syncthing
    eend $?
}

stop() {
    ebegin "Stopping syncthing"
    start-stop-daemon \
        --stop \
        --user "${SYNCTHING_USER}" \
        --group "${SYNCTHING_GROUP}" \
        --exec /usr/bin/syncthing
    eend $?
}
